# ops/systemd/cord-preview.service
[Unit]
Description=Cord Preview Docker Compose Application
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/cord-preview
Environment="HOME=/home/ubuntu"
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose stop
TimeoutStartSec=300

# Restart policy
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
